# docker/crontab/www-crontab
*/8 * * * * www /usr/local/bin/php /var/www/artisan app:parsing:polovniautomobili >> /var/log/cron.log 2>&1
*/1 * * * * www /usr/local/bin/php /var/www/artisan app:car-sender-command >> /var/log/cron.log 2>&1
