# docker/crontab/www-crontab
*/6 * * * * www /usr/local/bin/php /var/www/artisan app:parsing:polovniautomobili >> /var/log/cron.log 2>&1
