# docker/crontab/www-crontab
4 * * * * www /usr/local/bin/php /var/www/artisan app:hourly-car-command >> /var/log/cron.log 2>&1